# b.conf
start on started test/a
script 
  logger "value of foo is '$foo'"
end script
