# y.conf
start on eventa

script 
  echo "`date`: $UPSTART_JOB started" >>/tmp/eventa.log
  sleep 3
  echo "`date`: $UPSTART_JOB stoped" >>/tmp/eventa.log
end script
